# scripts/com.autoblog.smoke.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.autoblog.smoke</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/naseunghwan/Desktop/auto_blog_generator/scripts/smoke_publish_daily.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/naseunghwan/Desktop/auto_blog_generator</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>

    <!-- 매일 21:40 KST에 스모크 테스트를 실행한다. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>21</integer>
        <key>Minute</key>
        <integer>40</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>KeepAlive</key>
    <false/>

    <key>ProcessType</key>
    <string>Background</string>

    <key>StandardOutPath</key>
    <string>/Users/naseunghwan/Desktop/auto_blog_generator/logs/smoke.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/naseunghwan/Desktop/auto_blog_generator/logs/smoke.stderr.log</string>
</dict>
</plist>
